# fanfan-b-laravel/supervisor_soccer_live.conf
[program:SoccerWebSocket]
command=php artisan websocket:serve
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
directory=%(ENV_SOCCER_ROOT)s
user=root
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/supervisor/soccer/stdout_websocket.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_websocket.log ; 에러로그 정의

[program:SoccerScheduler]
command=php artisan schedule:work
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
directory=%(ENV_SOCCER_ROOT)s
user=root
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/supervisor/soccer/stdout_schedule.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_schedule.log ; 에러로그 정의


[program:SoccerLiveShortScheduler]
command=php artisan short-schedule:run
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
directory=%(ENV_SOCCER_ROOT)s
user=root
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/var/log/supervisor/soccer/stdout_short.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_short.log ; 에러로그 정의

[program:SoccerBaseQueue]
command=php artisan queue:work --queue=high,default,low,email,
process_name=%(program_name)s_%(process_num)02d
numprocs=4
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

# [program:SoccerLiveQueue]
# command=php artisan queue:work --queue=live_one,live_two,live_three,live_four,live_five,live_six,live_seven,live_eight,live_nine,live_ten
# process_name=%(program_name)s_%(process_num)02d
# numprocs=6
# directory=%(ENV_SOCCER_ROOT)s
# autostart=true
# autorestart=true
# user=root
# stopsignal=QUIT
# stopasgroup=true
# killasgroup=true
# redirect_stderr=true
# stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
# sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueOne]
command=php artisan queue:work --queue=live_one
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueTwo]
command=php artisan queue:work --queue=live_two
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueThree]
command=php artisan queue:work --queue=live_three
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueFour]
command=php artisan queue:work --queue=live_four
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueFive]
command=php artisan queue:work --queue=live_five
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueSix]
command=php artisan queue:work --queue=live_six
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueSeven]
command=php artisan queue:work --queue=live_seven
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueEight]
command=php artisan queue:work --queue=live_eight
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueNine]
command=php artisan queue:work --queue=live_nine
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SoccerLiveQueueTen]
command=php artisan queue:work --queue=live_ten
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의

[program:SimulationLiveQueue]
command=php artisan queue:work --queue=sim_seq
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=%(ENV_SOCCER_ROOT)s
autostart=true
autorestart=true
user=root
stopsignal=QUIT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/soccer/stdout_queue.log ; 출력로그 정의
sterr_logfile=/var/log/supervisor/soccer/sterr_queue.log ; 에러로그 정의